# lib/supervisor/conf.d/minecraft.conf
[program:minecraft]
startsecs=5
directory=/data/
command=java -Xmx1024M -Xms1024M -jar craftbukkit.jar
